# deploy/morpheus-line.service
[Unit]
Description=morpheus-line — Voice interface for Claude Code
After=network.target

[Service]
Type=simple
User=root
ExecStart=/usr/local/bin/morpheus-line
Environment=RUST_LOG=morpheus_line=info,tower_http=info
Restart=on-failure
RestartSec=5

[Install]
WantedBy=multi-user.target
